# claude-monitor.service
[Unit]
Description=Claude Monitor - Work Hour Tracking Daemon
Documentation=https://github.com/claude-monitor/system
After=network.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/.claude/monitor
ExecStart=/usr/local/bin/claude-monitor daemon --config /home/%i/.claude/monitor/config.json
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-monitor

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/%i/.claude
ProtectHome=tmpfs
BindReadOnlyPaths=/home/%i

# Resource limits
LimitNOFILE=8192
MemoryMax=512M
TasksMax=256

[Install]
WantedBy=default.target